# docs/sd_card/SpyderTallies.service
[Unit]
Description=Spyder Tally application by Derek Smithson (Knightware)
After=network.target

[Service]
WorkingDirectory=/opt/spyder-tally
ExecStart=/opt/spyder-tally/SpyderTallyControllerWebApp

# to query logs using journalctl, set a logical name here
SyslogIdentifier=SpyderTallies

# The install script will replace __USER__ with the installing user
User=__USER__

# Allow binding to port 80 without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Listen on port 80 directly (no reverse proxy needed)
Environment=ASPNETCORE_URLS=http://+:80

# ensure the service restarts after crashing
Restart=always
# amount of time to wait before restarting the service
RestartSec=5

[Install]
WantedBy=multi-user.target
